# conf/sidar.conf
description "sidar django application instance"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /home/sidar/django_projects/sidar/conf/start_server.sh 
